# dist/common/systemd/scylla-help.service.in
[Unit]
Description=Scylla Help

[Service]
Type=simple
EnvironmentFile=@@SYSCONFDIR@@/scylla-help.py
User=scylla
Group=scylla
ExecStart=/usr/lib/scylla/scylla-help.py -q version

[Install]
WantedBy=multi-user.target
