# dist/common/systemd/scylla-help.timer.in
[Unit]
Description=Run scylla help daily 
After=scylla-server.service
BindsTo=scylla-server.service

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
